# tmp102@.timer
[Unit]
Description=tmp102 timer %i

[Timer]
OnBootSec=1min
OnUnitActiveSec=10sec

[Install]
WantedBy=timers.target
