# tmp102.service
[Unit]
Description=tmp102

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/tmp102.js

[Install]
WantedBy=default.target
