# tmp102@.service
[Unit]
Description=tmp102 %i

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/tmp102.js %i

[Install]
WantedBy=default.target
